# diagnet/echo_server.service
[Unit]
Description=Synthetic TCP Echo Server
After=network-online.target
Wants=network-online.target
[Service]
Type=simple
EnvironmentFile=-/opt/css_diag_agent/diagnet.conf
ExecStart=/bin/bash -c '/usr/bin/env python /opt/css_diag_agent/diagnet/echo_server.py --bind 0.0.0.0 --port ${PEER_PORT:-9400}'
Restart=always
RestartSec=3s
User=root
[Install]
WantedBy=multi-user.target
